# data/data-provider-master-shortcut.socket
[Socket]
SocketUser=system
SocketGroup=system
ListenStream=/tmp/.shortcut.service
SocketMode=0666
SmackLabelIPIn=data-provider-master::shortcut
SmackLabelIPOut=@

Service=data-provider-master.service

[Unit]
Wants=data-provider-master.target
Before=data-provider-master.target

[Install]
WantedBy=sockets.target
