# data/data-provider-master-notification.socket
[Socket]
SocketUser=system
SocketGroup=system
ListenStream=/tmp/.notification.service
SocketMode=0666
SmackLabelIPIn=data-provider-master::notification
SmackLabelIPOut=@

Service=data-provider-master.service

[Unit]
Wants=data-provider-master.target
Before=data-provider-master.target

[Install]
WantedBy=sockets.target
