# data/data-provider-master.target
[Unit]
Description=data-provider-master sockets
DefaultDependencies=true
