# data/data-provider-master-utility.socket
[Socket]
SocketUser=system
SocketGroup=system
ListenStream=/tmp/.utility.service
SocketMode=0666
SmackLabelIPIn=data-provider-master::utility
SmackLabelIPOut=@

Service=data-provider-master.service

[Unit]
Wants=data-provider-master.target
Before=data-provider-master.target

[Install]
WantedBy=sockets.target
