# data/data-provider-master.path
[Unit]
Description=Check the WM_READY file to launch the data-provider-master

[Path]
PathExists=/run/.wm_ready
